# data-track.service
[Unit]
Description=YouTube data tracking service

[Service]
ExecStart=/home/udara/Documents/data_track/track.sh

[Install]
WantedBy=default.target
